# supervisor.conf
[supervisord]
nodaemon=true

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /app/uwsgi.ini
